# cron.weekly/rmoldvids
#!/bin/sh
find /var/stream/video -mtime +120 -delete
